# kali/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each mulituser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error
#
# In order to enable or disable this script just thange the execution
# bits.
#
# By default this script does nothing.

service ntp restart

#ipdown eth0
#ifup etho0

ifconfig eth0:0 192.168.2.100
route add -host 192.168.2.100 dev eth0:0

service apache restart
service mysql restart

service samba restart

exit 0
